# contrib/relay/orgii-mobile-relay.service
# ORGII Mobile Remote Relay — example systemd unit (system-level, hardened)
#
# Audience: self-hosters who want to run the relay as a long-lived
# service on a single-tenant Linux VPS, fronted by nginx (TLS termination).
#
# Companion deploy doc:
#   Documentation/MainApp/collaboration/mobile-remote-relay-deploy--0505.md
#
# ----------------------------------------------------------------------
# DIVERGENCE FROM `relay install` ON LINUX
# ----------------------------------------------------------------------
# `orgii-mobile-relay install` (see
# `src-tauri/crates/orgii-mobile-relay/src/cli/systemd.rs::render_unit`)
# generates a *user-level* unit at
# `~/.config/systemd/user/orgii-mobile-relay.service`, run via
# `systemctl --user`. That CLI-generated unit looks like this (verbatim
# from `render_unit`):
#
#   [Unit]
#   Description=ORGII Mobile Remote Relay
#   After=network-online.target
#   Wants=network-online.target
#
#   [Service]
#   Type=simple
#   ExecStart=/path/to/orgii-mobile-relay serve
#   Restart=on-failure
#   RestartSec=5
#
#   [Install]
#   WantedBy=default.target
#
# That unit is fine for a developer running the relay on their own
# desktop, but it does NOT bind a dedicated service user, has no
# sandboxing, and lives in the user's session scope. For an internet-
# reachable VPS you want the system-level, hardened unit below.
#
# This file IS NOT what `relay install` writes. Drop it manually at:
#   /etc/systemd/system/orgii-mobile-relay.service
#
# ----------------------------------------------------------------------
# Pre-reqs (one-time, as root):
#   useradd --system --no-create-home --shell /usr/sbin/nologin orgii-relay
#   install -d -o orgii-relay -g orgii-relay -m 0750 /var/lib/orgii-relay
#   install -m 0755 ./target/release/orgii-mobile-relay /usr/local/bin/
#
# Activation:
#   systemctl daemon-reload
#   systemctl enable --now orgii-mobile-relay
#   systemctl status orgii-mobile-relay
# ----------------------------------------------------------------------

[Unit]
Description=ORGII Mobile Remote Relay
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=orgii-relay
Group=orgii-relay
WorkingDirectory=/var/lib/orgii-relay
ExecStart=/usr/local/bin/orgii-mobile-relay serve --listen-addr 127.0.0.1:7878 --storage-path /var/lib/orgii-relay/orgii-relay.db
Restart=on-failure
RestartSec=2

# RUST_LOG override (optional). The relay falls back to its compiled
# default ("info") if RUST_LOG is unset; uncomment for verbose output.
# Environment=RUST_LOG=orgii_mobile_relay=debug,info

# Sandboxing. The relay only needs to read its binary, write to its
# state dir, and bind a loopback TCP port — nothing else.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/orgii-relay
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
